# debian/interchange.cron.daily
#! /bin/sh
#
# cron script to expire sessions
#
# Copyright (C) 2000 Stefan Hornburg <racke@linuxia.de>
   
/usr/lib/interchange/bin/expireall -f /etc/interchange/catalogs.cfg -r
