# debian/interchange.conffiles
/etc/cron.daily/interchange
/etc/interchange/catalogs.cfg
/etc/interchange/interchange.cfg
/etc/interchange/jobs/merchandising/topsellers
/etc/interchange/makecat.cfg
/etc/interchange/robots.cfg
/etc/init.d/interchange
/etc/logrotate.d/interchange
